# crontab
0 6 * * * docker exec r2-lens ruby main.rb >> /var/log/cron.log
