# crontab
0 6 * * * cd apps/r2-lens-original && docker compose run --rm r2-lens ruby main.rb >> /var/log/cron.log
